# wvwoverview.ini
[program:wvw-overview]
command=/usr/bin/python3.9 /home/sleider/WvWOverview/main.py
directory=/home/sleider/WvWOverview
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/home/sleider/logs/wvw-updater.log
stderr_logfile=/home/sleider/logs/wvw-updater.err
