# wvwoverview.ini
[program:wvw-overview]
command=/home/sleider/WvWOverview/WvWOverview
directory=/home/sleider/WvWOverview
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/home/sleider/logs/WvWOverview.log
stderr_logfile=/home/sleider/logs/WvWOverview.err
